# copy/opt/local/lib/svc/manifest/policyd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='policyd'>

  <service name='network/policyd' type='service' version='1'>
    <create_default_instance enabled='false' />
    <single_instance />

    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
       <service_fmri value='svc:/system/filesystem/local' />
    </dependency>

    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
       <service_fmri value='svc:/network/loopback' />
    </dependency>

    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
       <service_fmri value='svc:/spiped/percona:default' />
    </dependency>

    <exec_method type='method' name='start' exec='/opt/local/libexec/cbpolicyd' timeout_seconds='-1'>
       <method_context>
         <method_credential user='root' group='root' />
       </method_context>
    </exec_method>

    <exec_method type="method" name="stop"  exec=":kill" timeout_seconds="60"/>
  </service>

</service_bundle>
